# debian/postinst
#!/bin/sh
set -e
case "$1" in
    configure)
        # Create system user if it does not exist
        if ! id -u mp3-player > /dev/null 2>&1; then
            adduser --system --no-create-home --group \
                    --home /nonexistent \
                    --shell /usr/sbin/nologin \
                    mp3-player
        fi

        # Grant ALSA access
        if getent group audio > /dev/null 2>&1; then
            usermod -aG audio mp3-player
        fi

        # Lock down config
        chown root:mp3-player /etc/mp3-player/config.ini
        chmod 640 /etc/mp3-player/config.ini
        ;;
esac
#DEBHELPER#

echo ""
echo "mp3-player installed."
echo "  Edit config : /etc/mp3-player/config.ini"
echo "                Verify audio.device matches your sound card"
echo "                Run 'aplay -L' to list ALSA devices"
echo "  Enable      : systemctl enable --now mp3-player"
echo "  Play        : base64 -w0 song.mp3 | nc -q1 127.0.0.1 8566"
echo "  Play at vol : (printf \"volume=40\\n\"; base64 -w0 song.mp3) | nc -q1 127.0.0.1 8566"
echo ""
exit 0
